# ctf2019/suctf/homerouter/_homerouter.bin.extracted/squashfs-root/etc/uci-defaults/98_change_password
#!/bin/sh

system_set_password() {
    USER=$1
    PASS=$2
    (echo "$PASS"; sleep 1; echo "$PASS";) | passwd "$USER" 2>&1 >/dev/null
    return 0;
}

system_set_password root "nE7jA%5m"

